# etc/systemd/system/sysctl-config.service
[Unit]
Description=Use a modify sysctl config
After=default.target

[Service]
ExecStart=/sbin/sysctl -p

[Install]
WantedBy=default.target
